# dotfiles/linux/profile.d/kwin.sh
# Fix screen tearing in KDE
#
# See: https://www.reddit.com/r/kde/comments/99ykiw/screen_tearing_with_kde/
export KWIN_TRIPLE_BUFFER=1
