# dotfiles/linux/profile.d/nix.sh
if [ -e /home/jake/.nix-profile/etc/profile.d/nix.sh ]; then . /home/jake/.nix-profile/etc/profile.d/nix.sh; fi
